# systemd/wsvpn-client.service
[Unit]
Description=WSVPN Client - WebSocket VPN
Documentation=https://github.com/ntcjoel/wsvpn
After=network.target

[Service]
Type=simple
User=sq
Group=sq
WorkingDirectory=/home/sq/wsvpn-client
ExecStart=/home/sq/wsvpn-client/wsvpn-client
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wsvpn-client

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/sq/wsvpn-client
PrivateTmp=true

# Capabilities (instead of root)
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW

# Resource limits
LimitNOFILE=65536
LimitNPROC=1024

[Install]
WantedBy=multi-user.target
